# unison/dotfiles.cron
# MIN HOUR DAY MONTH WEEKDAY CMD

# Sync with AMD dotfiles every 5 minutes
  */5 *    *   *     *       if ifconfig tun0 &> /dev/null; then _log=~ks/.unison/dotfiles.cron.log.$(date --iso-8601); echo >> $_log; printf '=%.0s' {1..120} >> $_log; date "+\n[%T]" >> $_log; ~ks/bin/unison dotfiles >> $_log 2>&1; fi
