# unison/dotfiles.cron
# MIN HOUR DAY MONTH WEEKDAY CMD

# Sync with AMD dotfiles every 5 minutes
  */5 *    *   *     *       ifconfig tun0 &> /dev/null && unison dotfiles >| $HOME/.unison/dotfiles.cron.log.$(date --iso-8601) 2>&1
